# flask_turbo_boost/testdemo/deploy/supervisor.conf
[program:testdemo]
command=/var/www/testdemo/venv/bin/gunicorn -c deploy/gunicorn.conf wsgi:app
directory=/var/www/testdemo
user=root
autostart=true
autorestart=true
environment = MODE="PRODUCTION"